# home/private_dot_config/systemd/user/emacs.service.d/override.conf
# systemd user service には XDG_SESSION_TYPE が渡らず、/usr/bin/emacs ラッパーが
# X11 版にフォールバックするため、pgtk (Wayland) 版を明示する
[Service]
Environment=EMACS_TOOLKIT=wayland
