# home/private_dot_config/systemd/user/chezmoi-drift.service
[Unit]
Description=chezmoi drift check

[Service]
Type=oneshot
ExecStart=%h/.local/bin/chezmoi-drift-check
