# home/private_dot_config/systemd/user/chezmoi-drift.timer
[Unit]
Description=Daily chezmoi drift check

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
